# xci/source/info.apps/Manage-UIUC-V3/tags/tag-0.2-20200511/sbin/route_uiuc_daily_v3.crontab
0          5       * * * /bin/bash -l -c "/soft/warehouse-apps-1.0/Manage-UIUC-V3/PROD/sbin/route_uiuc_v3.rc start --ignore_dates"
